# supervisord.conf
[program:node]
command = node room/app.js
directory = /home/dotcloud/current
